# supervisor.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/etc/init.d/mysql start

[program:apache2]
command=/etc/init.d/apache2 start
